# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=warn

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ── PostgreSQL ────────────────────────────────────────────────────────────────
# priority=10 ensures postgres starts first; all other services depend on it.
# Runs as the postgres OS user. listen_addresses=* allows both localhost and
# external connections (port must still be exposed via -p in docker run).
[program:postgres]
command=/usr/lib/postgresql/14/bin/postgres
    -D /data/postgres
    -c port=%(ENV_PG_PORT)s
    -c listen_addresses=*
user=postgres
priority=10
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true

# ── Keycloak ──────────────────────────────────────────────────────────────────
# priority=20. Picks up KC_* env vars exported by entrypoint.sh into supervisord's
# environment. start-dev is appropriate for single-node / development deployments.
[program:keycloak]
command=/opt/keycloak/bin/kc.sh start-dev --import-realm
environment=
    KC_DB="postgres",
    KC_DB_URL="jdbc:postgresql://127.0.0.1:%(ENV_PG_PORT)s/keycloak",
    KC_DB_USERNAME="ndexserver",
    KC_DB_PASSWORD="%(ENV_KC_DB_PASSWORD)s",
    KC_BOOTSTRAP_ADMIN_USERNAME="%(ENV_KC_ADMIN_USER)s",
    KC_BOOTSTRAP_ADMIN_PASSWORD="%(ENV_KC_ADMIN_PASSWORD)s",
    KC_HOSTNAME_URL="%(ENV_KC_ISSUER_URL)s",
    KC_HOSTNAME_ADMIN_URL="%(ENV_KC_ISSUER_URL)s",
    KC_HTTP_PORT="%(ENV_KC_HTTP_PORT)s",
    KC_HTTP_MANAGEMENT_PORT="%(ENV_KC_MGMT_PORT)s",
    KC_HEALTH_ENABLED="true",
    KC_LOG_LEVEL="WARN"
priority=20
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true

# ── Solr ──────────────────────────────────────────────────────────────────────
# priority=20. Independent of Keycloak — starts in parallel.
# SOLR_HOME points to /data/solr/data where cores and configsets live.
# Runs as the solr OS user created in the Dockerfile.
[program:solr]
command=/opt/solr/bin/solr start -f -p %(ENV_SOLR_PORT)s -s /data/solr/data
environment=
    SOLR_HOME="/data/solr/data",
    SOLR_SECURITY_MANAGER_ENABLED="false"
user=solr
priority=20
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true

# ── MailHog ───────────────────────────────────────────────────────────────────
# priority=20. No dependencies — starts immediately.
[program:mailhog]
command=/usr/local/bin/mailhog
    -smtp-bind-addr 0.0.0.0:%(ENV_SMTP_PORT)s
    -api-bind-addr 0.0.0.0:%(ENV_MAILHOG_UI_PORT)s
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ── NDEx (Jetty) ──────────────────────────────────────────────────────────────
# priority=30: starts after all supporting services.
# autostart is controlled by NDEX_AUTOSTART:
#   false → devcontainer mode (developer runs 'mvn jetty:run' manually)
#   true  → docker run / Kubernetes mode (server starts automatically)
# In devcontainer mode the source is mounted at /workspaces/ndex-rest.
# For docker run / Kubernetes a pre-built artifact command should replace mvn jetty:run.
[program:ndex]
command=mvn jetty:run
directory=/workspaces/ndex-rest
environment=ndexConfigurationPath="/data/ndex/conf/ndex.properties"
autostart=%(ENV_NDEX_AUTOSTART)s
priority=30
autorestart=true
startretries=3
startsecs=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
