# docker/supervisord/keycloak.conf
# ── Keycloak ──────────────────────────────────────────────────────────────────
# All Keycloak config (DB credentials, ports, hostname, bootstrap admin) is in
# /opt/keycloak/conf/keycloak.conf — written by start.sh before supervisord starts.
[program:keycloak]
command=/opt/keycloak/bin/kc.sh start-dev --import-realm --log-level=info
environment=JAVA_TOOL_OPTIONS="-Dquarkus.log.level=INFO"
priority=20
autostart=true
autorestart=true
startretries=5
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopasgroup=true
